# deploy/quadlet/bindery.container
[Unit]
Description=Bindery Service
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/replace-me/bindery:latest
ContainerName=bindery
PublishPort=5670:5670
Volume=bindery-library.volume:/data:Z
Environment=BINDERY_LIBRARY_DIR=/data/library
Environment=BINDERY_TEMPLATE_DIR=/data/templates
EnvironmentFile=/etc/bindery/bindery.env
AutoUpdate=registry

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
